# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-vim-add-package/run
#!/usr/bin/with-contenv bash

if ! dpkg -s vim >/dev/null 2>&1; then
    echo "**** Adding vim to package install list ****"
    echo "vim" >> /mod-repo-packages-to-install.list
    echo "build-essential" >> /mod-repo-packages-to-install.list
    echo "wget" >> /mod-repo-packages-to-install.list
    echo "iputils-ping" >> /mod-repo-packages-to-install.list
    echo "nmap" >> /mod-repo-packages-to-install.list
    echo "telnet" >> /mod-repo-packages-to-install.list
    echo "podman" >> /mod-repo-packages-to-install.list
    echo "abc:100000:65536" >> /etc/subuid
else
    echo "**** vim already installed, skipping ****"
fi
